# etc-systemd-system/system/net-daily-summary.service.d/10-ntfy-env.conf
[Service]
EnvironmentFile=/etc/gandalf/ntfy.env
EnvironmentFile=-/etc/gandalf/secrets.env
