# etc-systemd-system/system/net-daily-summary.service.d/20-safe-wrapper.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/net-daily-summary-safe.sh
