# etc-systemd-system/system/net-daily-summary.service.d/50-only-if-bad.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/net-report.sh --daily-if-bad
