# etc-systemd-system/system/net-daily-summary.service.d/99-nonblocking.conf
[Service]
Restart=no
SuccessExitStatus=0 1 2 6 127
