# etc-systemd-system/system/net-health2.service.d/override.conf
# Gandalf observability override – non-blocking by design
[Service]
Restart=no
SuccessExitStatus=0 6 127
