# etc-systemd-system/system/gandalf-dns-selftest.timer
[Unit]
Description=Run Gandalf DNS self-test every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=1min
Unit=gandalf-dns-selftest.service

[Install]
WantedBy=timers.target
