# etc-systemd-system/system/gandalf-integrity.service
[Unit]
Description=Gandalf integrity check (protected files)
[Service]
Type=oneshot
ExecStart=/usr/local/bin/gandalf-integrity check
