# etc-systemd-system/system/gandalf-integrity-report.service
[Unit]
Description=Send weekly integrity report to ntfy
[Service]
Type=oneshot
ExecStart=/usr/local/bin/gandalf-integrity report
