# etc-systemd-system/system/gandalf-panel-smoke.service
[Unit]
Description=Gandalf Panel smoke test (status+auth+ddns)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gandalf-panel-smoke
