# etc-systemd-system/system/gandalf-panel-cert-rotate.service
[Unit]
Description=Rotate Gandalf panel TLS cert and reload nginx
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/gandalf-panel-cert-rotate
