# etc-systemd-system/system/net-report-daily.service
[Unit]
Description=Send daily merged network report to ntfy
[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-report.sh --daily
