# etc-systemd-system/system/net-health.service
[Unit]
Description=Periodic network health check (fix via recover-net if needed)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-health.sh
