# etc-systemd-system/system/net-daily-summary.timer
[Unit]
Description=Daily 09:00 network summary
[Timer]
OnCalendar=09:00
Persistent=true
Unit=net-daily-summary.service
[Install]
WantedBy=timers.target
