# etc-systemd-system/system/net-monthly-report.service
[Unit]
Description=Send monthly Gandalf network summary to ntfy
[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-monthly-report.sh
