# etc-systemd-system/system/net-health2.service
[Unit]
Description=Network health check (auto-recover if needed)
[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-report.sh --health
