# etc-systemd-system/system/net-daily-summary.service
[Unit]
Description=Send daily network summary to ntfy
[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-daily-summary.sh
