# etc-systemd-system/system/net-health.timer
[Unit]
Description=Run network health check every 2 minutes

[Timer]
OnBootSec=1m
OnUnitActiveSec=2m
RandomizedDelaySec=15s
Unit=net-health.service

[Install]
WantedBy=timers.target
