# etc-systemd-system/system/pihole-blockall-off.timer
[Unit]
Description=Disable blockall daily (kids/guest)

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
