# systemd/cloudflared.service
[Unit]
Description=cloudflared DNS over HTTPS proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/cloudflared --no-autoupdate --config /etc/cloudflared/config.yml proxy-dns
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
