# etc-systemd-system/system/gandalf-dns-selftest.service.d/override.conf
# Gandalf observability override – non-blocking by design
[Service]
Restart=no
SuccessExitStatus=0 1 6 127
