# boot-kit/firstrun.service
[Unit]
Description=Gandalf Zero-Touch First Run
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash /boot/firmware/boot-kit/firstrun.sh
StandardOutput=append:/var/log/gandalf-firstrun.log
StandardError=append:/var/log/gandalf-firstrun.log
Restart=no

[Install]
WantedBy=multi-user.target